# docker/supervisord.conf
# All-in-one process supervisor for the Unraid / single-container deployment.
#
# AlbuMine needs three processes: Redis (the ARQ queue broker), the ARQ worker
# (image processing) and the FastAPI web app. For a one-click Unraid install
# they run together in one container, managed by supervisord. The docker-compose
# dev setup instead overrides the container command to run each as its own
# service, so this file is only used by the default (all-in-one) command.
#
# supervisord itself runs as root (so it can forward child output to the
# container's stdout/stderr) but drops every managed program to the
# unprivileged 'albumine' user — created by the entrypoint from PUID/PGID.
# Writable runtime files go to /tmp.

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

# XML-RPC endpoint for the web UI's "restart server" button. Bound to the
# container-local loopback only — nothing outside the container can reach it.
[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:redis]
user=albumine
command=redis-server --save "" --appendonly no --bind 127.0.0.1 --port 6379
priority=10
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
user=albumine
command=albumine-worker
priority=20
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
user=albumine
command=albumine
priority=30
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
